# core/scripts/podman/jarvis.container.in
# JARVIS gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into jarvis's ~/.config/containers/systemd/
# {{JARVIS_HOME}} is replaced at install time.

[Unit]
Description=JARVIS gateway (rootless Podman)

[Container]
Image=jarvis:local
ContainerName=jarvis
UserNS=keep-id
Volume={{JARVIS_HOME}}/.jarvis:/home/node/.jarvis
EnvironmentFile={{JARVIS_HOME}}/.jarvis/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
